# pkg/systemd/wallpaper-service.service
[Unit]
Description=Desktop Wallpaper (video/image/color)
After=graphical-session.target monitor-detect.service
Wants=monitor-detect.service

[Service]
Type=simple
WorkingDirectory=%h/Systems/desktop/wallpaper-service
ExecStart=%h/Systems/desktop/wallpaper-service/src/wallpaper_service/wallpaper_service.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
