# pkg/systemd/wallpaper-service.service
[Unit]
Description=Desktop Wallpaper (video/image/color)
After=graphical-session.target monitor-detect.service
Wants=monitor-detect.service

[Service]
Type=simple
WorkingDirectory=-/home/nick/.local/share/wallpaper-service
ExecStart=/usr/local/bin/wallpaper-service
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
